# backup-daemon.sh
#! bin/bash
### BEGIN INIT INFO
# Short-Description: Code Backup and Deployment Daemon
# Description: Daemon that backs up code and deploys changed code.
### END INIT INFO
DAEMON=/home/markb/Documents/Systems\ Software/Assignment\ 1/backupDaemon

exit 0
